# threadgpio.service
[Unit]
Description=Thread GPIO 

[Service]
ExecStart=/home/pi/unisal/atividades/gpioBtnC/gpioBtnThread 

[Install]
WantedBy=multi-user.target
